# backend/cron/cleanup.cron
0 1 * * * /usr/local/bin/python /app/script/cleanup_reservations.py >> /var/log/cleanup.log 2>&1