# backend/cron/cleanup.cron
0 1 * * * root /app/script/run_cleanup.sh >> /app/logs/run_cleanup.log 2>&1
